# ansible_collections/tesselate/k8s/roles/k8s_common/files/openrc/kubelet.initd
#!/sbin/openrc-run

description='The core Kubernetes node agent'
supervisor=supervise-daemon

command='/usr/bin/kubelet'
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

: ${output_log:=/var/log/${RC_SVCNAME}/stdout.log}
: ${error_log:=/var/log/${RC_SVCNAME}/stderr.log}

depend() {
	need net
}
